# roles/speedtest/templates/speedtest-run.timer.j2
[Unit]
Description=Hourly speedtest

[Timer]
OnCalendar=hourly
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
